# ops/eventarb.service
[Unit]
Description=EventArb Trading Bot
After=network.target
Requires=network.target

[Service]
Type=simple
User=eventarb
Group=eventarb
WorkingDirectory=/opt/eventarb-bot
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/eventarb-bot
ExecStart=/usr/bin/python3 /opt/eventarb-bot/app.py
Restart=always
RestartSec=5

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eventarb-bot

[Install]
WantedBy=multi-user.target

